# misc/health-bot.crontab
MAILTO=""
*/2 * * * * /bin/bash /home/hackermans/src/health-bot/misc/cronwrapper.sh
0 0 * * * /usr/sbin/logrotate -s $HOME/logs/logrotate.state misc/health-bot.logrotate
